# os/unix/linux/tech/booting/systemd/integrating_lsb_scripts
*** LSB scripts will not necessarily run after all other systemd services ***
*** For maximum control, replace scripts with native systemd services ***

To integrate a legacy LSB boot script:

Create the script in /etc/rc.d/init.d and make it executable
The script must have the following entries in comments:

# chkconfig: 2345 99 01
# description: LSB integration test script
#
### BEGIN INIT INFO
# Provides: lsbtest
# Required-Start: $local_fs $network $remote_fs $named $time
# Required-Stop: $local_fs $network $remote_fs $named
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: LSB test
# Description: LSB test
### END INIT INFO

where the "2345" is the list of runlevels in which the script will be started
(could be any set of digits between 0 and 6 inclusive)
the "99" is the sequence number for the start script (S99)
the "01" is the sequence number for the stop script (K01)
Other pairs of numbers besides 99/01 can be used, as needed

Run "chkconfig scriptname reset" to create S or K symlinks as appropriate
   in ALL /etc/rc.d/rc[0-6].d directories
S scripts are created in dirs listed in the first parameter to the commented
    chkconfig line in the script
K scripts are created in all other rc[0-6].d dirs

Running "systemctl enable scriptname" is equivalent to "chkconfig scriptname on"
   which automatically creates S scripts in the /etc/rc.d/rc[2345].d dirs
   but does not honor the parameter to the commented chkconfig line in the
   script and does not create any K scripts
